# bau-os.service
[Unit]
Description=Bau-OS KI-Assistent (Telegram Bot)
Documentation=https://docs.bau-os.at
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=bauos
WorkingDirectory=/opt/bau-os
EnvironmentFile=/opt/bau-os/.env
Environment=LANG=de_AT.UTF-8 LC_ALL=de_AT.UTF-8
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bau-os

# Resource Limits
MemoryMax=512M
MemoryHigh=384M
CPUQuota=100%
TasksMax=64

# Sicherheit
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/bau-os/logs /opt/bau-os/.chat_id /opt/bau-os-vault /opt/bau-os/data /opt/bau-os/tools
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
